# deploy/turnable-manager-panel.service
[Unit]
Description=Turnable Manager Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=turnable-manager
Group=turnable-manager
Environment=TURNABLE_MANAGER_ADMIN_PASSWORD=change-this-first-run-only
ExecStart=/usr/local/bin/turnable-manager-panel -config /etc/turnable-manager-panel/config.json
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/turnable-manager-panel /etc/turnable

[Install]
WantedBy=multi-user.target
